# assets/supervisor.bind.conf
[program:bind9]
autostart = true
command = /usr/sbin/named -g
redirect_stderr = true
stdout_logfile = /var/log/bind/bind.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10
